# system/daemons/gofed-resources-client.service
[Unit]
Description=Resource client GC
After=multi-user.target

[Service]
Type=idle
Environment=PYTHONPATH=/home/jchaloup/Projects/gofed
ExecStart=/usr/bin/python /home/jchaloup/Projects/gofed/infra/system/daemons/resourceclientdaemon.py

[Install]
WantedBy=multi-user.target
